# odysseus-ui.service
# systemd unit for the web UI. Prefer ./install-service.sh, which fills in the
# bracketed tokens below from your environment (user, install dir, bind host/port,
# venv) and refuses to install if any token is left unresolved.
#
# To install by hand instead: replace each token and copy to
# /etc/systemd/system/odysseus-ui.service.
[Unit]
Description=Odysseus UI
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__INSTALL_DIR__
# Bind to loopback by default (front with a reverse proxy or `tailscale serve`);
# never 0.0.0.0 on a privileged box.
ExecStart=__INSTALL_DIR__/__VENV__/bin/uvicorn app:app --host __BIND_HOST__ --port __APP_PORT__
Restart=always
RestartSec=3
EnvironmentFile=-__INSTALL_DIR__/.env

[Install]
WantedBy=multi-user.target
